# q-gateway.service
[Unit]
Description=Q Gateway Service
After=network.target
Wants=network-online.target

[Service]
Type=exec
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/huixin/aiops_v2

# Environment
Environment=PATH=/home/ubuntu/huixin/aiops_v2/.venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/ubuntu/huixin/aiops_v2
Environment=SOP_DIR=/home/ubuntu/huixin/aiops_v2/sop
Environment=TASK_DOC_PATH=/home/ubuntu/huixin/aiops_v2/task_instructions.md
Environment=QTTY_HOST=127.0.0.1
Environment=QTTY_PORT=7682
Environment=QTTY_PING=55
Environment=QTTY_TIMEOUT=20
Environment=ALERT_JSON_PRETTY=1

# Service execution
ExecStartPre=/bin/bash -c 'pkill -f "uvicorn.*app:app" || true; pkill -f "ttyd.*q_entry.sh" || true; sleep 2'
ExecStart=/home/ubuntu/huixin/aiops_v2/gateway/start.sh
ExecStop=/bin/bash -c 'pkill -f "uvicorn.*app:app" || true; pkill -f "ttyd.*q_entry.sh" || true'
ExecReload=/bin/kill -HUP $MAINPID

# Process management
Restart=always
RestartSec=10
KillMode=mixed
TimeoutStartSec=60
TimeoutStopSec=30
KillSignal=SIGTERM

# Resource limits
MemoryMax=4G
MemoryHigh=3G
CPUQuota=300%
TasksMax=200

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/ubuntu/huixin/aiops_v2
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=q-gateway

[Install]
WantedBy=multi-user.target
