# systemd/aiops-ttyd.service
[Unit]
Description=AIOps TTYD Q Chat Service
After=network.target

[Service]
Type=forking
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/huixin/aiops_v2
ExecStart=/home/ubuntu/huixin/aiops_v2/scripts/start_ttyd.sh
ExecStop=/bin/bash -c 'if [[ -f /home/ubuntu/huixin/aiops_v2/logs/ttyd.pid ]]; then kill $(cat /home/ubuntu/huixin/aiops_v2/logs/ttyd.pid); fi'
PIDFile=/home/ubuntu/huixin/aiops_v2/logs/ttyd.pid
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
